# packages/linux/.config/systemd/user/pulseaudio-dlna.service
[Unit]
Description=PulseAudio DLNA Service

[Install]
WantedBy=default.target

[Service]
Type=simple
PrivateTmp=true
ExecStart=/usr/bin/pulseaudio-dlna
